# ansible/roles/backend/templates/backend.service.j2
[Unit]
Description=FastAPI backend
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/backend
Environment="DATABASE_URL=postgresql://{{ rds_user }}:{{ rds_password }}@{{ rds_endpoint }}:5432/{{ rds_db }}"
ExecStart=/opt/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
